# udev/99-sentri-auto-wifi-config.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z]1", RUN+="/usr/local/bin/sentri-auto-wifi-config.sh"
